# deploy/docker/cp-dav/sync/crontab
* * * * * root flock -w 1 /var/run/sync-nfs.lock bash ${SYNC_HOME}/sync-nfs.sh
